# systemd/.config/systemd/user/freya-sessions.service
[Unit]
Description=Create tmux sessions for bazaar and garden
# Bind to the graphical session so WAYLAND_DISPLAY is in the systemd user
# environment before the tmux server is created. niri --session imports
# WAYLAND_DISPLAY and then reaches graphical-session.target; starting earlier
# (e.g. default.target) freezes an empty WAYLAND_DISPLAY into the tmux server.
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=PATH=%h/.bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/.bin/,freya-sessions.sh

[Install]
WantedBy=graphical-session.target
